# units/virtiofsd@.service.in
#  SPDX-License-Identifier: LGPL-2.1-or-later
#
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=virtiofsd Service for directory %I
Requires=virtiofsd@%i.socket
After=virtiofsd@%i.socket

[Service]
Type=notify
ExecStart={{VIRTIOFSD_PATH}} --fd ${LISTEN_FDS} --shared-dir %I
RuntimeDirectory=virtiofsd/%i
RuntimeDirectoryPreserve=yes
Sockets=virtiofsd@%i.socket
